# appengine/flexible/websockets/additional-supervisord.conf
[program:socket-server]
command = php %(ENV_APP_DIR)s/socket-server.php
enviroment = PORT="8000"
autostart = true
autorestart = true
priority = 10
stopwaitsecs = 20

