# appengine/flexible/websockets/additional-supervisord.conf
[program:socket-server]
command = php %(ENV_APP_DIR)s/socket-server.php
enviroment = PORT="8000"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes=0
user = root
autostart = true
autorestart = true
priority = 10
stopwaitsecs = 20

